# docker/conf/supervisor/supervisord.ini
[supervisord]
nodaemon=true
pidfile=/var/run/supervisord.pid
logfile=/var/www/storage/supervisor/supervisord.log

[program:php-fpm]
command=/usr/local/sbin/php-fpm -c /usr/local/etc/php/php.ini -y /usr/local/etc/php-fpm.conf -F

[program:redis]
command=redis-server
directory=/var/www/storage/supervisor
autostart=true
autorestart=true
stdout_logfile=/var/www/storage/supervisor/redis.log
stderr_logfile=/var/www/storage/supervisor/redis-error.log

[program:worker-default]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/artisan queue:work --queue=default --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/storage/supervisor/queue-default.log
stopwaitsecs=3600